# deploy/db_writer.service
# deploy/db_writer.service
# Install: sudo cp deploy/db_writer.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable db_writer
#          sudo systemctl start db_writer

[Unit]
Description=Maverick Telemetry — DB Writer
Documentation=https://github.com/AlexTs-dev/maverick-telemetry-hub
# Must start after Mosquitto is fully up
After=network.target mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/maverick-telemetry-hub
ExecStart=/home/pi/maverick-telemetry-hub/venv/bin/python db_writer.py

# Restart policy — always restart on failure, with 5s delay
Restart=always
RestartSec=5

# Give it 10s to shut down cleanly before killing it
TimeoutStopSec=10

# Log to journald — view with: journalctl -u db_writer -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=db_writer

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=MAVERICK_DB_PATH=/home/pi/maverick_telemetry.db

[Install]
WantedBy=multi-user.target
