# deploy/express_bridge.service
# deploy/express_bridge.service
# Install: sudo cp deploy/express_bridge.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable express_bridge
#          sudo systemctl start express_bridge

[Unit]
Description=Maverick Telemetry — Express Bridge
Documentation=https://github.com/AlexTs-dev/maverick-telemetry-hub
After=network.target mosquitto.service db_writer.service
Requires=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/maverick-telemetry-hub/server

# Ensure the React build exists before starting.
# Run `npm run build` in client/ after every frontend change.
ExecStartPre=/bin/test -f /home/pi/maverick-telemetry-hub/client/dist/index.html

ExecStart=/usr/bin/node index.js

Restart=always
RestartSec=5
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=express_bridge

Environment=NODE_ENV=production
Environment=MAVERICK_DB_PATH=/home/pi/maverick_telemetry.db
# Create /home/pi/maverick-telemetry-hub/server/.env with:
# ANTHROPIC_API_KEY=your_key_here
# PORT=3000
EnvironmentFile=-/home/pi/maverick-telemetry-hub/server/.env

[Install]
WantedBy=multi-user.target
